# config/systemd_service_files/bluetooth.service
[Unit]
Description=Bluetooth Listener 
After=network.target

[Service]
ExecStart=/home/mdp23/MDP_RPI/.venv/bin/python /home/mdp23/MDP_RPI/bluetooth/bluetooth_listener.py
Restart=always
User=pi
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
