# config/systemd_service_files/pi_cam.service
[Unit]
Description=Picamera2 Camera Service
After=network.target

[Service]
Type=simple
User=mdp23
Group=mdp23
WorkingDirectory=/home/mdp23/MDP_RPI
# Use system python; update the path to your script accordingly
ExecStart=/usr/bin/python pi_camera/camera_service.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
