# config/systemd_service_files/stm_comm.service
[Unit]
Description=STM Communication Service
After=network.target

[Service]
WorkingDirectory=/home/mdp23/MDP_RPI
ExecStart=/home/mdp23/MDP_RPI/.venv/bin/python -m stm_comm.serial_comm
Restart=always
User=pi
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/home/mdp23/MDP_RPI

[Install]
WantedBy=multi-user.target
